# deploy/kdeskdash.service
# systemd unit for the kdeskdash dashboard appliance.
#
# Install on the Pi (rpidash2):
#   sudo install -D -m644 deploy/kdeskdash.service /etc/systemd/system/kdeskdash.service
#   sudo install -D -m600 deploy/kdeskdash.env.example /etc/kdeskdash/kdeskdash.env  # then edit
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now kdeskdash
#
# The binary is deployed to /usr/local/bin/kdeskdash by `cmake --build build-pi --target deploy`.
[Unit]
Description=kdeskdash touch dashboard
# Redis is optional; start after it when present so remote control/persistence
# are available immediately, but do not fail if it is absent.
After=redis-server.service
Wants=redis-server.service
# Never permanently give up restarting the appliance. With systemd's default
# start limit (StartLimitBurst=5 / StartLimitIntervalSec=10s) a transient
# early-boot crash loop would wedge the unit in `failed` until someone SSHes in
# — unacceptable for a headless boot-to-dashboard panel. 0 disables the limit.
# (This key lives in [Unit] on systemd >= 229, which covers Trixie.)
StartLimitIntervalSec=0

[Service]
# DRM master + evdev access require root on this panel.
User=root
ExecStart=/usr/local/bin/kdeskdash
# Leading '-' makes the env file optional: every KDESKDASH_* var has a built-in
# default, so a missing file must not block the appliance from starting.
EnvironmentFile=-/etc/kdeskdash/kdeskdash.env
Restart=always
RestartSec=2

# --- Sandboxing -----------------------------------------------------------
# Shrink the blast radius of a root process while preserving the two devices the
# panel must reach: the DRM node (/dev/dri/*) and the touch evdev
# (/dev/input/*). We deliberately do NOT set PrivateDevices= / DevicePolicy=,
# which would hide those nodes; the filesystem hardening below leaves /dev
# untouched, so DRM master and touch keep working.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
# Only writable state: the icons-mode favourites file. StateDirectory creates
# /var/lib/kdeskdash and makes just that path writable under ProtectSystem=strict
# (matching KDESKDASH_ICONS_FAVORITES's default). Screenshots land in
# PrivateTmp's /tmp; the env file and runtime font are read-only.
StateDirectory=kdeskdash

[Install]
WantedBy=multi-user.target
